# box_configuration/jetson/roscore.service
[Unit]
Description=ROScore service
After=network-online.target

[Service]
Type=forking
User=rsl
ExecStart=/bin/sh -c ". /opt/ros/your_ros_distro/setup.sh; . /etc/ros/env.sh; roscore & while ! echo exit | nc 192.168.2.51 11311 > /dev/null; do sleep 1; done"

[Install]
WantedBy=multi-user.target